# btsp.service
[Unit]
Description=Bluetooth service

[Service]
Type=forking
ExecStart=/etc/init.d/btsp

