# btsp.service
[Unit]
Description=Bluetooth service

[Service]
Type=forking
ExecStart=/usr/bin/btsp

[Install]
WantedBy=multi-user.target

